# systemd/certdeployer.timer
[Unit]
Description=Run CertDeployer daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target